# docker/wordpress/supervisord.conf
[supervisord]

[program:mysqld]
command=/usr/bin/mysqld_safe

[program:httpd]
command=/etc/apache2/foreground.sh
stopsignal=6

[program:vsftpd]
command=vsftpd
